# docker_assets/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=service php8.1-fpm restart

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
